# root/etc/s6-overlay/s6-rc.d/service-proton/finish
#!/command/with-contenv bash
# shellcheck shell=bash

rm -rf "${CONFIG_DIR}/forwarded_port"
