# root/etc/s6-overlay/s6-rc.d/service-proton/run
#!/command/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

vpn_remote=$(grep -P -o -m 1 '(?<=^Endpoint)(\s{0,})[^:]+' < "${CONFIG_DIR}/wireguard/${VPN_CONF}.conf"| sed -e 's~^[=\ ]*~~')

while true; do
    if ! natpmpc -g "${vpn_remote}"; then
        echo "[WARNING] ProtonVPN endpoint \"${vpn_remote}\" does not support port forwarding!"
    else
        port=$(natpmpc -g "${vpn_remote}" -a 1 0 udp 60 | grep -P -o -m 1 '(?<=Mapped public port\s)\d+')
        if [[ -n "${port}" ]]; then
            natpmpc -g "${vpn_remote}" -a 1 0 tcp 60
            webui_https="http"
            if grep -q 'WebUI\\HTTPS\\Enabled=true' "${CONFIG_DIR}/config/qBittorrent.conf"; then
                webui_https="https"
            fi
            curl -fsSL --insecure -X POST -d "json={\"random_port\": false}" "${webui_https}://localhost:${WEBUI_PORTS%%/*}/api/v2/app/setPreferences"
            curl -fsSL --insecure -X POST -d "json={\"listen_port\": ${port}}" "${webui_https}://localhost:${WEBUI_PORTS%%/*}/api/v2/app/setPreferences"
            echo "[INFO] ProtonVPN forwarded port is \"${port}\"."
        else
            echo "[WARNING] ProtonVPN port forwarding failed!"
        fi
    fi
    sleep 45
done
