# root/etc/s6-overlay/s6-rc.d/init-perms/up
find /etc/s6-overlay/s6-rc.d -name "run*" -execdir chmod +x {} +
